# cursor/ai-consultant/deploy/lvt-ai-consultant.service
[Unit]
Description=LVT AI consultant orchestrator (OpenAI + Bitrix tools)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/www-root/data/www/lvtgroup.ru/cursor/ai-consultant
EnvironmentFile=/var/www/www-root/data/www/lvtgroup.ru/cursor/ai-consultant/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
